# contrib/init/bitcoind.conf
description "Bitcoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env polcoind_BIN="/usr/bin/polcoind"
env polcoind_USER="bitcoin"
env polcoind_GROUP="bitcoin"
env polcoind_PIDDIR="/var/run/polcoind"
# upstart can't handle variables constructed with other variables
env polcoind_PIDFILE="/var/run/polcoind/polcoind.pid"
env polcoind_CONFIGFILE="/etc/bitcoin/bitcoin.conf"
env polcoind_DATADIR="/var/lib/polcoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # polcoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$polcoind_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run polcoind."
        echo "The setting must appear in $polcoind_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Bitcoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$polcoind_PIDDIR"
    chmod 0755 "$polcoind_PIDDIR"
    chown $polcoind_USER:$polcoind_GROUP "$polcoind_PIDDIR"
    chown $polcoind_USER:$polcoind_GROUP "$polcoind_CONFIGFILE"
    chmod 0660 "$polcoind_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$polcoind_PIDFILE" \
    --chuid $polcoind_USER:$polcoind_GROUP \
    --exec "$polcoind_BIN" \
    -- \
    -pid="$polcoind_PIDFILE" \
    -conf="$polcoind_CONFIGFILE" \
    -datadir="$polcoind_DATADIR" \
    -disablewallet \
    -daemon

